# contrib/systemd/control.target
[Unit]
Description=Control System
After=network.target
